# root/etc/init/doraemon.conf
# ================= DO NOT MODIFY THIS FILE =================
#
# Manual changes will be lost when this file is regenerated.
#
# Please read the developer's guide, which is available
# at https://dev.nethesis.it/projects/nethserver/wiki/NethServer
# original work from http://www.contribs.org/development/
#
# Copyright (C) 2013 Nethesis S.r.l.
# http://www.nethesis.it - support@nethesis.it
#
#
# 10base doraemon - Apache instance that runs the doraemon service
#

description "Apache instance that runs the doraemon service"
author "Paolo Asperti <paolo@asperti.com>"
respawn
respawn limit 10 5
expect fork
start on started network
stop on stopping network


#
# 20pre_start
#
pre-start script
    STAT=`/sbin/e-smith/config getprop doraemon status`
    if test "x$STAT" = "xenabled" \
       -a -f "/etc/e-smith/db/configuration/defaults/doraemon/status" \
       -a -x /usr/sbin/doraemon; then
       exit 0;
    fi
    stop; exit 0;
end script


#
# 30script
#
exec /usr/sbin/doraemon -f /etc/httpd/doraemon/httpd.conf
